# root/etc/uci-defaults/99-openclaw
#!/bin/sh
# luci-app-openclaw — 首次安装初始化脚本

# 已移除创建 openclaw 系统用户的逻辑，以适应自定义安装目录通常为 NTFS/FAT32 的外置 USB 驱动器（外置驱动器不支持 chown 且往往默认为 root 用户权限）。
# 这里保留注释，旧版升级上来的用户其残留 openclaw 账号会在卸载时自动清理。

# 动态获取安装路径
OC_INSTALL_DIR=$(uci -q get openclaw.main.install_dir 2>/dev/null || echo "/opt/openclaw")
[ -z "$OC_INSTALL_DIR" ] && OC_INSTALL_DIR="/opt/openclaw"

# 创建数据目录
mkdir -p "$OC_INSTALL_DIR/data/.openclaw"
mkdir -p "$OC_INSTALL_DIR/node"
mkdir -p "$OC_INSTALL_DIR/global"
chown -R root:root "$OC_INSTALL_DIR" 2>/dev/null || true

# 生成随机 Token (如果尚未设置)
CURRENT_TOKEN=$(uci -q get openclaw.main.token)
if [ -z "$CURRENT_TOKEN" ]; then
	TOKEN=$(head -c 24 /dev/urandom | hexdump -e '24/1 "%02x"' 2>/dev/null || dd if=/dev/urandom bs=24 count=1 2>/dev/null | od -An -tx1 | tr -d ' \n' | head -c 48)
	uci set openclaw.main.token="$TOKEN"
	uci commit openclaw
fi

# 生成 PTY Token (如果尚未设置)
CURRENT_PTY_TOKEN=$(uci -q get openclaw.main.pty_token)
if [ -z "$CURRENT_PTY_TOKEN" ]; then
	PTY_TOKEN=$(head -c 24 /dev/urandom | hexdump -e '24/1 "%02x"' 2>/dev/null || dd if=/dev/urandom bs=24 count=1 2>/dev/null | od -An -tx1 | tr -d ' \n' | head -c 48)
	uci set openclaw.main.pty_token="$PTY_TOKEN"
	uci commit openclaw
fi

exit 0
